# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; PHP-FPM — handles all HTTP → FastCGI requests from nginx
[program:php-fpm]
command=php-fpm -F -R
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; Nginx — HTTP server, proxies to PHP-FPM on 127.0.0.1:9000
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
; Wait for PHP-FPM to be ready before starting nginx
startsecs=3

; Generation worker — polls Replicate API for personalization jobs
[program:generation-worker]
command=php /srv/sylius/bin/console app:personalization:process-generation-jobs --loop --limit=10 --sleep-seconds=2 --max-runtime=31536000 --env=prod
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
directory=/srv/sylius
priority=30
; Wait for app to be fully ready before starting worker
startsecs=10

; Messenger worker — processes Stripe webhooks and other async jobs with retry/failure transport
[program:messenger-worker]
command=php /srv/sylius/bin/console messenger:consume async --time-limit=3600 --memory-limit=256M --env=prod -vv
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
directory=/srv/sylius
priority=40
startsecs=10
